# service.startup/jalousie.service
[Unit]
Description=jalousie
After=pigpiod.service

[Service]
User=<user>
ExecStart=<full path to>/jalousie.js

[Install]
WantedBy=multi-user.target
